# examples/remote-worker/supervisord.conf
; supervisord configuration for the remote-runtime example.
;
; It launches one runtime (the standalone control plane, owning the SQLite data
; store) and two worker hosts that connect to it over WebTransport.
;
; Build the binaries first (see README.md):
;   go build -o bin/runtime ../../cmd/runtime
;   go build -o bin/worker .
;
; Then run:
;   supervisord -c supervisord.conf
;
; Note: the workers reconnect to the runtime with backoff, so the startup order
; below (runtime first via "priority") is best-effort, not strict — a worker
; started before the runtime is ready will keep retrying until it connects.

[supervisord]
logfile=supervisord.log
pidfile=supervisord.pid
nodaemon=true

[program:runtime]
command=%(here)s/bin/runtime -config %(here)s/config.yaml
directory=%(here)s
priority=1
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker1]
command=%(here)s/bin/worker -worker-address 127.0.0.1:8081 -actor-host-address 127.0.0.1:7571 -runtime-address 127.0.0.1:7400
directory=%(here)s
priority=10
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker2]
command=%(here)s/bin/worker -worker-address 127.0.0.1:8082 -actor-host-address 127.0.0.1:7572 -runtime-address 127.0.0.1:7400
directory=%(here)s
priority=10
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
