# examples/worker/supervisord.conf
; supervisord configuration for the local-runtime example.
;
; It launches two worker hosts. Each embeds its own actor host and SQLite data
; store; the workers discover and invoke each other peer-to-peer.
;
; Build the binary first (see README.md):
;   go build -o bin/worker .
;
; Then run:
;   mkdir -p worker1 worker2
;   supervisord -c supervisord.conf
;
; Each worker creates a data.db relative to its working directory, so they need
; separate subdirectories to avoid sharing the same SQLite file.

[supervisord]
logfile=supervisord.log
pidfile=supervisord.pid
nodaemon=true

[program:worker1]
command=%(here)s/bin/worker -worker-address 127.0.0.1:8081 -actor-host-address 127.0.0.1:7571
directory=%(here)s/worker1
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:worker2]
command=%(here)s/bin/worker -worker-address 127.0.0.1:8082 -actor-host-address 127.0.0.1:7572
directory=%(here)s/worker2
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
